# app-backup/bacula/files/2.4.0/bacula-all-init
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-backup/bacula/files/2.4.0/bacula-all-init,v 1.1 2008/06/07 15:41:29 wschlich Exp $

depend() {
	need net
	if [ "${LOCAL_DB}" -eq 1 ]; then
		need "%database%"
	fi
	use dns
}

start() {
	local retval
	declare -i retval=0
	for s in ${SERVICES} ; do
		case ${s} in
			sd)
				ebegin "Starting bacula storage daemon"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-sd \
					-- ${SD_OPTIONS}
				let retval+=$?
				;;
			fd)
				ebegin "Starting bacula file daemon"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-fd \
					-- ${FD_OPTIONS}
				let retval+=$?
				;;
			dir)
				ebegin "Starting bacula director"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-dir \
					-- ${DIR_OPTIONS}
				let retval+=$?
				;;		
			*)
				eerror "Unrecognized service found in SERVICES: '${s}'"
				;;
		esac
	done
	eend ${retval}
}

stop() {
	local retval
	declare -i retval=0
	for s in ${SERVICES} ; do
		case ${s} in
			sd)
				ebegin "Stopping bacula storage daemon"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-sd.*.pid 
				let retval+=$?
				;;
			fd)
				ebegin "Stopping bacula file daemon"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-fd.*.pid
				let retval+=$?
				;;
			dir)
				ebegin "Stopping bacula director"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-dir.*.pid
				let retval+=$?
				;;		
			*)
				eerror "Unrecognized service found in SERVICES: '${s}'"
				;;
		esac
	done
	eend ${retval}
}
